# service/backup_files.service
[Unit]
Description=Copy specified files to google drive using gio. Useful to make backups

[Service]
ExecStart=%h/.bin/backup_files.bash
